# packaging/systemd/fbforward-geoip-update.service
[Unit]
Description=Update fbforward GeoIP databases
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/fbforward-geoip
ExecStart=/usr/local/libexec/fbforward/update-geoip.sh
