# packaging/systemd/fbforward-geoip-update.timer
[Unit]
Description=Periodic fbforward GeoIP database update

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
