# deploy/packaging/debian/postinst
#!/usr/bin/env sh
set -eu

SERVICE_USER="fbforward"
SERVICE_GROUP="fbforward"
CONFIG_PATH="/etc/fbforward/config.yaml"

if ! getent group "$SERVICE_GROUP" >/dev/null 2>&1; then
  groupadd --system "$SERVICE_GROUP"
fi

if ! id -u "$SERVICE_USER" >/dev/null 2>&1; then
  useradd --system --gid "$SERVICE_GROUP" --home /nonexistent --shell /usr/sbin/nologin "$SERVICE_USER" || \
    useradd --system --gid "$SERVICE_GROUP" --home /nonexistent --shell /sbin/nologin "$SERVICE_USER"
fi

if [ -f "$CONFIG_PATH" ]; then
  chown root:"$SERVICE_GROUP" "$CONFIG_PATH"
  chmod 0640 "$CONFIG_PATH"
fi

systemctl daemon-reload || true
systemctl enable --now fbforward.service || true
