# deploy/packaging/debian/prerm
#!/usr/bin/env sh
set -eu

systemctl disable --now fbforward.service || true
systemctl daemon-reload || true
